# etc/rusted.timer
[Unit]
Description=timer for the rusted config backup service

[Timer]
OnStartupSec=1m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
